# router-config/lede-bak/etc/crontabs/root
*     *       *       *       *   echo `date` > /kuaicdn/app/ipes/var/log/ipes/cron.log && sh /kuaicdn/app/ipes/bin/check start >> /kuaicdn/app/ipes/var/log/ipes/cron.log  2>&1
03 4 * * * sleep 5 && touch /etc/banner && reboot   //每天4点03分路由器自动重启
0 4 * * * /etc/coremark.sh
